# systemd/alertmanager-discord-bridge.service
[Unit]
Description=Alertmanager Discord Bridge
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu
ExecStart=/usr/bin/python3 /home/ubuntu/alertmanager-discord-bridge.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
